# playbook/templates/blink.service.j2
[Unit]
Description=Gewitter blink 

[Service]
Type=simple
ExecStart=/opt/blink.sh
ExecStartPre=/opt/blink.sh init
Restart=always

[Install]
WantedBy=multi-user.target
